# deploy/powerbi-gateway.service
[Unit]
Description=Power BI Embedded Gateway
After=network-online.target postgresql-15.service
Wants=network-online.target
Requires=postgresql-15.service

[Service]
Type=simple
User=test01
WorkingDirectory=/home/test01/projects/powerbi-embedded
EnvironmentFile=/home/test01/projects/powerbi-embedded/.env
ExecStartPre=/usr/bin/mkdir -p /home/test01/projects/powerbi-embedded/logs
ExecStartPre=/home/test01/.pyenv/versions/3.11.9/bin/python3 /home/test01/projects/powerbi-embedded/scripts/migrate_report_meta.py
ExecStart=/home/test01/.pyenv/versions/3.11.9/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8247
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=powerbi-gateway

[Install]
WantedBy=multi-user.target
